# deploy/systemd/cita-bot.service
[Unit]
Description=CitaConsular slot monitor Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/bot/app
ExecStart=/root/bot/app/venv/bin/python -m cita_bot
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
